# ansible/playbooks/roles/systemd/mount/templates/systemd.mount.j2
[Unit]
Description={{ mnt.description }}

[Mount]
What={{ mnt.what }}
Where={{ mnt.where }}
Type={{ mnt.type }}
Options={{ mnt.options }}

[Install]
WantedBy=multi-user.target
